# deploy/external-kms-laptop/usb-presence/kerplace-kms-seal.service
[Unit]
Description=KerPlace custody: lazy-umount s3fs + seal the KMS (custody USB removed)
# Oneshot, triggered by the udev rule and by kerplace-kms-presence.
[Service]
Type=oneshot
ExecStart=/usr/local/sbin/kerplace-kms-seal
