# deploy/external-kms-laptop/systemd-user/kerplace-usb-watch.service
# ~/.config/systemd/user/kerplace-usb-watch.service
# Custody USB insertion agent: on INSERT it runs 'kpc enable' (brings up the
# tunnel, unseals with a graphical prompt, mounts the buckets); on REMOVAL it
# stops the tunnel(s). The instantaneous seal+unmount is done by the system
# service (udev + kerplace-kms-presence). Edge-triggered (no .path unit loop).
# It runs in the user session -> inherits DISPLAY/DBUS -> kdialog works.
#
# Install:  systemctl --user enable --now kerplace-usb-watch.service
# (and linger:  sudo loginctl enable-linger $USER  — so it starts after a boot)
[Unit]
Description=KerPlace: custody USB agent (enable on insert / stop tunnel on removal)

[Service]
ExecStart=/usr/local/bin/kpc watch
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
