# deploy/external-kms-laptop/systemd-user/kerplace-tunnel.service.example
# ~/.config/systemd/user/kerplace-<HOST>-tunnel.service
# Custody tunnel (laptop -> KerPlace host): a single SSH connection carries
#   -R 8200  custody (the host reaches your local KMS)
#   -L 9000  data    (you reach the host's S3, encrypted in transit by SSH)
# The post-quantum KEX is pinned by policy. It does NOT auto-start at boot: its
# lifecycle follows the USB's presence (kpc enable starts it, kpc watch stops it).
# Replace <HOST>, <PORT>, <USER>, and the KEX to match the host's sshd:
#   ML-KEM (OpenSSH >= 9.9):  mlkem768x25519-sha256
#   PQC fallback (>= 8.5):    sntrup761x25519-sha512
[Unit]
Description=KerPlace custody tunnel (laptop -> <HOST>)
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/usr/bin/ssh -N -T -p <PORT> \
  -R 127.0.0.1:8200:127.0.0.1:8200 \
  -L 127.0.0.1:9000:127.0.0.1:9000 \
  -o KexAlgorithms=sntrup761x25519-sha512 \
  -o StrictHostKeyChecking=accept-new \
  -o BatchMode=yes -o ConnectTimeout=10 \
  -o ServerAliveInterval=15 -o ServerAliveCountMax=3 \
  -o ExitOnForwardFailure=yes \
  <USER>@<HOST>
Restart=on-failure
RestartSec=5
