# deploy/external-kms-laptop/systemd/kerplace-kms-tunnel.service
[Unit]
Description=KerPlace custody-mode SSH tunnel (custody -R 8200 + data -L 9000)
Documentation=file:docs/OFFHOST_KMS_CUSTODY.md
Wants=network-online.target
After=network-online.target

[Service]
# A single SSH connection carries both directions of custody mode:
#   -R 127.0.0.1:8200 : CUSTODY — the remote host reaches the local Vault (KMS).
#   -L 127.0.0.1:9000 : DATA    — this laptop reaches KerPlace's S3, encrypted in
#                                 transit by the SSH channel (ML-KEM).
# No tunnel means no keys, no data and no endpoint: fail-closed throughout.
#
# Hardened: the post-quantum KEX is PINNED by policy (if the server does not offer
# ML-KEM the connection fails rather than downgrading) and the host key is pinned
# against a dedicated known_hosts (no TOFU). ExitOnForwardFailure=yes makes a
# forward that cannot bind tear the connection down (fail-closed). %h = the
# user's home.
#
# ADJUST: YOUR_USER@YOUR_HOST and the key / known_hosts paths for your deployment.
# The server-side key must be restricted in authorized_keys:
#   restrict,port-forwarding,permitlisten="127.0.0.1:8200",permitopen="127.0.0.1:9000" ssh-ed25519 AAAA... custody
# Note: in a later phase the key will come from an ssh-agent loaded off the USB,
# not from disk.
ExecStart=/usr/bin/ssh -N -T \
  -R 127.0.0.1:8200:127.0.0.1:8200 \
  -L 127.0.0.1:9000:127.0.0.1:9000 \
  -i %h/.ssh/kerplace-custody.key \
  -o KexAlgorithms=mlkem768x25519-sha256 \
  -o StrictHostKeyChecking=yes \
  -o UserKnownHostsFile=%h/.config/kerplace/known_hosts \
  -o BatchMode=yes \
  -o ConnectTimeout=10 \
  -o ServerAliveInterval=15 -o ServerAliveCountMax=3 \
  -o ExitOnForwardFailure=yes \
  YOUR_USER@YOUR_HOST
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
