# deploy/external-kms-laptop/systemd/kerplace-kms-tunnel.service
[Unit]
Description=Túnel SSH del modo custodia KerPlace (custodia -R 8200 + datos -L 9000)
Documentation=file:docs/OFFHOST_KMS_CUSTODY.md
Wants=network-online.target
After=network-online.target

[Service]
# Una sola conexión SSH lleva las dos direcciones del modo custodia:
#   -R 127.0.0.1:8200 : CUSTODIA — el host remoto alcanza el Vault local (KMS).
#   -L 127.0.0.1:9000 : DATOS    — este portátil alcanza el S3 de KerPlace (T1),
#                                  cifrado en tránsito por el canal SSH (ML-KEM).
# Sin túnel no hay claves, ni datos, ni endpoint: fail-closed total.
#
# Endurecido (T4): KEX post-cuántico FIJADO por política (si el servidor no ofrece
# ML-KEM, la conexión falla, no degrada) + host key pinning contra un known_hosts
# dedicado (sin TOFU). ExitOnForwardFailure=yes hace que un forward que no prende
# tumbe la conexión (fail-closed). %h = home del usuario.
#
# AJUSTA: YOUR_USER@YOUR_HOST y las rutas de clave / known_hosts a tu despliegue.
# La clave del servidor debe estar restringida en authorized_keys:
#   restrict,port-forwarding,permitlisten="127.0.0.1:8200",permitopen="127.0.0.1:9000" ssh-ed25519 AAAA... custodia
# Nota (fase 2): la clave vendrá de ssh-agent cargada desde el USB, no de disco.
ExecStart=/usr/bin/ssh -N -T \
  -R 127.0.0.1:8200:127.0.0.1:8200 \
  -L 127.0.0.1:9000:127.0.0.1:9000 \
  -i %h/.ssh/kerplace-custody.key \
  -o KexAlgorithms=mlkem768x25519-sha256 \
  -o StrictHostKeyChecking=yes \
  -o UserKnownHostsFile=%h/.config/kerplace/known_hosts \
  -o BatchMode=yes \
  -o ConnectTimeout=10 \
  -o ServerAliveInterval=15 -o ServerAliveCountMax=3 \
  -o ExitOnForwardFailure=yes \
  YOUR_USER@YOUR_HOST
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
